# service_template.ini
[program:service_template]
command = /python3/bin/uwsgi.ini --ini /service_template/uwsgi-sit.ini
directory = /service_template
startsecs = 10
startsetries = 5
autostart = true
autorestart = true
stdout_logfile = /service_template/logs/stdout.log
stdout_logfile_maxbytes = 10MB
redirect_stderr = true
user = root
stopasgroup = true
killasgroup = true
